# deploy/eftforge.service
[Unit]
Description=EFTForge FastAPI Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/eftforge/backend
Environment="PATH=/opt/eftforge/backend/venv/bin"
EnvironmentFile=/opt/eftforge/backend/.env
ExecStart=/opt/eftforge/backend/venv/bin/python -m gunicorn main:app \
    -w 1 \
    -k uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --access-logfile - \
    --error-logfile -
Restart=always
RestartSec=5
MemoryMax=1500M
RuntimeMaxSec=28800

[Install]
WantedBy=multi-user.target
